# 18.06/package/base-files/files/etc/hotplug.d/phy/00-phy
#!/bin/sh

logger -t phy "$ACTION $IFNAME"
